# launchd/com.nanoclaw-codex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nanoclaw-codex</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{PROJECT_ROOT}}/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{PROJECT_ROOT}}</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{PROJECT_ROOT}}/logs/nanoclaw-codex.log</string>
  <key>StandardErrorPath</key>
  <string>{{PROJECT_ROOT}}/logs/nanoclaw-codex.error.log</string>
</dict>
</plist>
